# lab-rip-hw/net/r4/etc/rc.local
#!/bin/sh -e

iptables -F
iptables -F -t nat
iptables -X

# ethX connects router to "Internet" via NAT

#iptables -t nat -A POSTROUTING -o ethX -j MASQUERADE

exit 0
